# deb/postinst
#!/bin/sh

systemctl enable stratos
systemctl start stratos
