# deb/prerm
#!/bin/sh

systemctl stop stratos
systemctl disable stratos
